# internal-flash/etc/rc.local

/sd/connectedclients/decrypt_leases.sh
#Add your commands above this
exit 0
