# internal-flash/etc/rc.local

/sd/connectedclients/decrypt_leases.sh
#Add your commands above this
mkdir -p /sd/logntouch /sd/log/dnsmasq.log
exit 0
